# systemd/livekit.service
[Unit]
Description=LiveKit SFU (media for live sessions)
After=network.target

[Service]
Type=simple
User=mentra
ExecStart=/usr/local/bin/livekit-server --config /etc/mentra/livekit.yaml
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
